# n8n/rootfs/etc/s6-overlay/s6-rc.d/n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# Runs N8N script
# ==============================================================================

bashio::log.info "Starting N8N..."

export N8N_USER_FOLDER="/config/n8n/data"
export N8N_CONFIG_FILES=/config/n8n/config.json

if bashio::config.has_value 'public_url'; then
  export N8N_EDITOR_BASE_URL="$(bashio::config 'public_url')"
fi

exec n8n start
