# n8n/rootfs/etc/s6-overlay/s6-rc.d/n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# Runs N8N script
# ==============================================================================

bashio::log.info "Starting N8N..."

export N8N_PATH="$(bashio::addon.entry)"

if bashio::config.exists 'log_level'; then
  export N8N_LOG_LEVEL="$(bashio::string.lower "$(bashio::config 'log_level')")"
fi

exec n8n start
