# n8n/rootfs/etc/services.d/n8n/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables

## Get the 'message' key from the user config options.


## Print the message the user supplied, defaults to "Hello World..."


## Run your program

#!/bin/sh
if [ "$#" -gt 0 ]; then
  # Got started with arguments
  node "$@"
else
  # Got started without arguments
  n8n
fi