# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# N8N init script, runs before any other service
# ==============================================================================

bashio::log.info "Oh, hi! I'm the init script!"

export VUE_APP_URL_BASE_API="127.0.0.1:$(bashio::addon.ingress_port)"