# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# N8N init script, runs before any other service
# ==============================================================================

bashio::log.info "Oh, hi! I'm the init script!"

ADDRESS="$(bashio::addon.ip_address)"
PORT="$(bashio::addon.ingress_port)"

bashio::log.info "Configuring N8n on address: ${ADDRESS}:${PORT}"

export VUE_APP_URL_BASE_API="127.0.0.1:"