# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# N8N init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring N8N..."

declare port
declare protocol

# Ensure configuration data directory exists
if ! bashio::fs.directory_exists '/config/data'; then
    mkdir -p /config/data || bashio::exit.nok "Failed to create n8n configuration directory"
fi

if ! bashio::config.true 'ssl'; then
  protocol="http";
else
  protocol="https";
fi

port="$(bashio::addon.port 3111)"
if ! bashio::var.has_value "${port}"; then
  port=3111
fi

bashio::var.json \
    protocol "$protocol" \
    host "$(bashio::addon.hostname)" \
    port "^$port" \
    ssl_cert "/ssl/$(bashio::config 'certfile')" \
    ssl_key "/ssl/$(bashio::config 'keyfile')" \
    db_type "$(bashio::config 'db_type')" \
    db_name "$(bashio::config 'db_name')" \
    db_host "$(bashio::config 'db_host')" \
    db_port "^$(bashio::config 'db_port')" \
    db_user "$(bashio::config 'db_user')" \
    db_password "$(bashio::config 'db_password')" \
    log_level "$(bashio::string.lower "$(bashio::config 'log_level')")" \
    | tempio \
        -template /etc/n8n/config.gtpl \
        -out /data/config.json

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /data/config.json)"
fi

