# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# N8N init script, runs before any other service
# ==============================================================================

bashio::log.info "Configuring N8N..."

# Ensure configuration exists
if ! bashio::fs.directory_exists '/config/n8n/'; then
    mkdir -p /config/n8n/data || bashio::exit.nok "Failed to create n8n configuration directory"

    # Copy in template files
    rm /config/n8n/.config.json
    cp /etc/n8n/config.json /config/n8n/
fi

if bashio::config.exists 'public_url'; then
  sed -i "s/%%PUBLIC_URL%%/"$(bashio::config 'public_url')"/" "/config/n8n/config.json"
else
  sed -i "s/%%PUBLIC_URL%%//" "/config/n8n/config.json"
fi

# Log level
if bashio::config.exists 'log_level'; then
  sed -i "s/%%LOG_LEVEL%%/$(bashio::string.lower "$(bashio::config 'log_level')")/" "/config/n8n/config.json"
fi

if [ "$(bashio::config 'log_level')" == "debug" ]; then
  bashio::log.info "$(cat /config/n8n/config.json)"
fi

