# n8n/rootfs/etc/s6-overlay/s6-rc.d/init-n8n/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: N8N
# N8N init script, runs before any other service
# ==============================================================================

bashio::log.info "Oh, hi! I'm the init script!"

export N8N_HOST=127.0.0.1
export N8N_PORT=28431

if bashio::config.exists 'log_level'; then
  export N8N_LOG_LEVEL="$(bashio::string.lower "$(bashio::config 'log_level')")"
fi