# n8n/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: N8N
# Runs the Nginx daemon
# ==============================================================================

# Wait for N8N to become available
bashio::net.wait_for 46543

bashio::log.info "Starting NGinx..."
exec nginx
